# docker/services/fpm/run
#!/bin/bash
exec /etc/init.d/php5-fpm start